# server/supervisor/local.ini
[program:harvest_inc-local]
directory = /home/devuser/.virtualenvs/harvest_inc-env/harvest_inc
command = /home/devuser/.virtualenvs/harvest_inc-env/bin/uwsgi --ini /home/devuser/.virtualenvs/harvest_inc-env/harvest_inc/server/uwsgi/local.ini
autostart = true
autorestart = true
